# deploy/nq_executor.service
# systemd unit for nq_executor (NQ ORB live trading C++ daemon)
#
# Install:  sudo cp deploy/nq_executor.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable nq_executor
#           sudo systemctl start nq_executor
#
# Logs:     journalctl -u nq_executor -f
# Stop:     sudo systemctl stop nq_executor   (sends SIGTERM → emergency flatten)
# Status:   sudo systemctl status nq_executor

[Unit]
Description=NQ ORB Live Executor (C++)
Documentation=https://github.com/neoptolemospapadiofantous/rithmic_engine
After=network-online.target postgresql-16.service
Wants=network-online.target
Requires=postgresql-16.service

[Service]
Type=simple

User=opc
Group=opc
WorkingDirectory=/home/opc/rithmic_engine

# Load credentials from .env (KEY=value format, one per line)
EnvironmentFile=/home/opc/rithmic_engine/.env

# Hard-block start if the NO_DEPLOY lockfile is present.
# Remove this file manually after verifying config for the session.
ExecStartPre=/bin/bash -c \
    '[ ! -f /home/opc/rithmic_engine/NO_DEPLOY ] || \
     { echo "NO_DEPLOY lockfile present — refusing to start executor"; exit 1; }'

ExecStart=/home/opc/rithmic_engine/build/nq_executor \
    /home/opc/rithmic_engine/config/live_config.json

# Allow up to 90 s for Rithmic login + market data subscription before marking failed
TimeoutStartSec=90s

# Restart automatically on crash; NOT on clean exit (trade day end)
Restart=on-failure
RestartSec=30s
# Max 3 restarts in 5 minutes — avoids hammering Rithmic on auth failure
StartLimitIntervalSec=300s
StartLimitBurst=3

# Allow up to 30 s on SIGTERM for position flatten + DB session flush
TimeoutStopSec=30s

# SIGTERM triggers the executor's emergency flatten / graceful shutdown handler
KillMode=process
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL
TimeoutKillSec=60s

# Route stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nq_executor

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
